# eye_api/tools/supervisor-eye.ini
[program:eye-api]
command = bash /data/eye/eye_api/tools/start-api.sh
autostart = true
stdout_logfile = /data/eye/eye_api/logs/api.log
redirect_stderr = true

[program:eye-ws]
command = bash /data/eye/eye_api/tools/start-ws.sh
autostart = true
stdout_logfile = /data/eye/eye_api/logs/ws.log
redirect_stderr = true

[program:eye-worker]
command = bash /data/eye/eye_api/tools/start-worker.sh
autostart = true
stdout_logfile = /data/eye/eye_api/logs/worker.log
redirect_stderr = true

[program:eye-monitor]
command = bash /data/eye/eye_api/tools/start-monitor.sh
autostart = true
stdout_logfile = /data/eye/eye_api/logs/monitor.log
redirect_stderr = true

[program:eye-scheduler]
command = bash /data/eye/eye_api/tools/start-scheduler.sh
autostart = true
stdout_logfile = /data/eye/eye_api/logs/scheduler.log
redirect_stderr = true